# cli/lnk-clib/t/examples/echo.service.in
[Unit]
Description=Echo Server

[Service]
Type=simple
WorkingDirectory=$CWD
ExecStart=$CWD/bin/echo
