# cli/lnk-clib/t/examples/echo.socket
[Unit]
Description=Socket Activated Echo Server

[Socket]
ListenStream=/tmp/echo.sock
FileDescriptorName=echo
