# bins/example-systemd-config/lnk-gitd.service
[Unit]
Description=lnk git daemon

[Service]
ExecStart=<your lnk-gitd> --linger-timeout 10000 --push-seeds --fetch-seeds
Environment=RUST_LOG=info
