# bins/example-systemd-config/lnk-gitd.socket
[Socket]
ListenStream=127.0.0.1:9987
FileDescriptorName=ssh

[Install]
WantedBy=sockets.target
